# packages/qmp-nycmesh/files/etc/hotplug.d/iface/95-nodog
#!/bin/sh
if [ "$ACTION" = ifup ]; then
    if [ "$DEVICE" = "br-lan" ]; then
        #/etc/init.d/nodogsplash enabled 
	/etc/init.d/nodogsplash stop
	sleep 120
	/etc/init.d/nodogsplash start
    fi
fi
