# packages/qmp-nycmesh/files/etc/init.d/znycmesh_autoconf
#!/bin/sh /etc/rc.common                                                                                                                          
                                                                                                                                                  
START=99                                                                                                                                          
STOP=99                                                                                                                                           
QMP_CONTROL_FILE="/qmp_configured"                                                                                                                
NYCMESH_CONTROL_FILE="/nycmesh_configured"                                                                                                        
                                                                                                                                                  
boot() {                                                                                                                                          
        echo boot                                                                                                                                 
        if [ -e "$QMP_CONTROL_FILE" ]; then                                                                                                       
                echo "qMp configured, running NYC Mesh conf"                                                                                      
                configure                                                                                                                         
        else                                                                                                                                      
                echo "qMp is not configured yet"                                                                                                  
        fi                                                                                                                                        
}                                                                                                                                                 
                                                                                                                                                  
configure() {                                                                                                                                     
        echo "entering configure"                                                                                                                 
        if [ ! -e "$NYCMESH_CONTROL_FILE" ]; then                                                                                                 
        echo "Configuring NYC Mesh"                                                                                                               
        initial_conf                                                                                                                              
        else                                                                                                                                      
                echo "NYC Mesh is already configured"                                                                                             
		start_up
        fi                                                                                                                                        
}                                                                                                                                                 
                                                                                                                                                  
initial_conf() {                                                                                                                                  
        echo "running initial conf"                                                                                                               
        /etc/nycmesh/nycmesh_initial_conf.sh                                                                                                              
        touch $NYCMESH_CONTROL_FILE                                                                                                               
# runs via nycmesh_conf.sh
#	qmpcontrol configure_wifi
	reboot
}
start_up() {
	echo Startup
	tincd -n nycmesh
}
                                                                                                                                                  
start() {                                                                                                                                         
        echo start                                                                                                                                
                                                                                                                                                  
}                                                                                                                                                 
                                                                                                                                                  
stop() {                                                                                                                                          
        echo stop                                                                                                                                 
}                 
