# packages/qmp-nycmesh/files/etc/uci-defaults/znycmesh_banner
#!/bin/sh
cat /etc/nycmesh/nycmesh.banner >> /etc/banner
exit 0
