# packages/qmp-nycmesh/files/etc/uci-defaults/znycmesh_tinc
#!/bin/sh
/etc/nycmesh/tinc_conf.sh
exit 0
